# machines/noti/claude-lid-watcher.service
[Unit]
Description=Claude Code: inhibidor de tapa + watcher de suspensión
Documentation=https://github.com/TonioAmc/dotfiles-claude

[Service]
Type=simple
# El servicio SOSTIENE el inhibidor (systemd-inhibit envuelve al watcher): si el
# watcher muere, el inhibidor se suelta solo y systemd lo reinicia. El watcher
# decide cuándo Suspend() según tapa + sesiones activas + piso de batería.
ExecStart=/usr/bin/systemd-inhibit --what=handle-lid-switch --who=claude --why="Claude Code monitoreando lid" --mode=block /home/antolin/.local/bin/claude-lid-watcher.sh
# exit 0 = "no hay claude" -> no reiniciar (queda inactivo hasta el próximo
# UserPromptSubmit). Crash/OOM/señal = failure -> reiniciar.
Restart=on-failure
RestartSec=2

# Sin [Install]: NO se habilita al boot. Lo arranca claude-inhibit-start.sh
# bajo demanda (systemctl --user start) cuando Claude empieza un turno.
